# scripts/com.ricksi.workbench-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  每日自动同步数据到仓库（本地数据源只在这台 Mac 上，所以用 launchd 而非云端 cron）
  安装：
    cp scripts/com.ricksi.workbench-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.ricksi.workbench-sync.plist
  卸载：
    launchctl unload ~/Library/LaunchAgents/com.ricksi.workbench-sync.plist
  注意：WorkingDirectory 如果移动了仓库位置需要同步修改。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ricksi.workbench-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>sync</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/rick/Claude_Code/Rick's Personal/rick-homepage</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd 环境没有交互 shell 的 PATH 与代理 alias：
         PATH 让 npm/node/git/claude 可寻址；代理变量供 collect-frontier 抓外网与 claude CLI 使用 -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HTTPS_PROXY</key>
    <string>http://127.0.0.1:7897</string>
    <key>HTTP_PROXY</key>
    <string>http://127.0.0.1:7897</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/workbench-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/workbench-sync.err</string>
</dict>
</plist>
