# scripts/com.ricksi.workbench-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  每日自动同步数据到仓库（本地数据源只在这台 Mac 上，所以用 launchd 而非云端 cron）
  安装：
    cp scripts/com.ricksi.workbench-sync.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.ricksi.workbench-sync.plist
  卸载：
    launchctl bootout gui/$(id -u)/com.ricksi.workbench-sync
  注意：WorkingDirectory 固定指向 main 专用生产 worktree，不能改回日常开发 checkout。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ricksi.workbench-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>sync</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/rick/Claude_Code/Rick's Personal/rick-homepage-sync</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd 环境没有交互 shell 的 PATH 与代理 alias：
         PATH 让 npm/node/git/gh 可寻址；Codex 使用配置中的 App 内绝对路径；代理供抓取与 GitHub 使用 -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HTTPS_PROXY</key>
    <string>http://127.0.0.1:7897</string>
    <key>HTTP_PROXY</key>
    <string>http://127.0.0.1:7897</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <!-- 持久日志（不在 /tmp，重启不丢）。简要「每次运行结果」见 ~/Library/Logs/workbench-sync.log（由 sync.mjs 追加）。 -->
  <key>StandardOutPath</key>
  <string>/Users/rick/Library/Logs/workbench-sync.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rick/Library/Logs/workbench-sync.stderr.log</string>
</dict>
</plist>
